# docker/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid
user=root

[program:php-fpm]
command=php-fpm -F -O
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/php-fpm.err
stdout_logfile=/var/log/supervisor/php-fpm.out

[program:nginx]
command=nginx -g "daemon off;"
autostart=true
autorestart=true
stderr_logfile=/var/log/supervisor/nginx.err
stdout_logfile=/var/log/supervisor/nginx.out

[program:queue-worker]
command=php /var/www/html/artisan queue:work --sleep=3 --tries=3 --timeout=120
autostart=true
autorestart=true
numprocs=1
stderr_logfile=/var/log/supervisor/queue.err
stdout_logfile=/var/log/supervisor/queue.out
user=www-data

[program:vite]
command=npm run dev
directory=/var/www/html
autostart=true
autorestart=true
user=www-data
stdout_logfile=/var/log/supervisor/vite.log
stderr_logfile=/var/log/supervisor/vite-error.log
